# zensical/rootfs/etc/s6-overlay/s6-rc.d/nginx/finish
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# nginx finish handler.
#
# s6-supervise restarts the service ONLY when finish exits 0.
# - Signal kill / exit 0 → halt the container (manual stop or graceful shutdown).
# - Non-zero exit → exit 0 (restart nginx).  nginx rarely crashes on its own;
#   a config error would cause an obvious crash-loop.  We don't halt the
#   container so the zensical-watcher keeps running and the user can read logs.
# ==============================================================================
# shellcheck disable=SC2155
readonly exit_code_container=$(</run/s6-linux-init-container-results/exitcode)
readonly exit_code_service="${1}"
readonly exit_code_signal="${2}"

bashio::log.info \
  "nginx exited with code ${exit_code_service}" \
  "(by signal ${exit_code_signal})"

if [[ "${exit_code_signal}" -ne 0 ]]; then
  # Killed by signal — transient (manual stop / resource limit).  Halt.
  if [[ "${exit_code_container}" -eq 0 ]]; then
    echo $((128 + exit_code_signal)) > /run/s6-linux-init-container-results/exitcode
  fi
  exec /run/s6/basedir/bin/halt
elif [[ "${exit_code_service}" -eq 0 ]]; then
  # Clean exit — transient (graceful shutdown).  Halt.
  exec /run/s6/basedir/bin/halt
else
  # Non-zero exit — restart nginx via exit 0 (s6-supervise restarts when
  # finish exits 0).  nginx rarely crashes on its own; when it does, a
  # config error is likely → crash-loop would be obvious and the user can
  # stop the add-on.  We don't halt the container so the watcher keeps
  # running and the user can read logs.
  bashio::log.warning \
    "nginx: crashed with exit ${exit_code_service} — restarting"
  exit 0
fi
