# zensical/rootfs/etc/s6-overlay/s6-rc.d/nginx/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Serve the built site at /var/lib/zensical/site/ on the ingress port.
# init-zensical guarantees the directory exists before this longrun starts.
# ==============================================================================
mkdir -p /run/nginx
exec nginx -g "daemon off;"
