# nemotron-asr-c/rootfs/etc/s6-overlay/s6-rc.d/nemotron-asr-c/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ============================================================================
# s6 longrun for nemotron-asr-c — Wyoming STT server
# ============================================================================
set -euo pipefail

cd /usr/src/app || exit 1

flags=()
flags+=(
    --uri "tcp://0.0.0.0:10370"
    --model "$(bashio::config 'model')"
    --quantization "$(bashio::config 'quantization')"
    --chunk-size "$(bashio::config 'chunk_size')"
)

bashio::config.has_value 'hf_token' && flags+=('--hf-token' "$(bashio::config 'hf_token')")
bashio::config.has_value 'hotword_boost' && flags+=('--hotword-boost' "$(bashio::config 'hotword_boost')")

# Hotwords are a YAML list; join into one newline-separated --hotwords value.
hotwords=$(jq -r '(.hotwords // []) | join("\n")' /data/options.json 2>/dev/null || echo "")
[ -n "${hotwords}" ] && flags+=('--hotwords' "${hotwords}")

bashio::config.true 'speech_enhancement' && flags+=(--speech-enhancement)
bashio::config.true 'debug_logging' && flags+=(--debug)

exec python3 -m wyoming_nemotron_asr_c "${flags[@]}"
