# nemo-asr-c/rootfs/etc/s6-overlay/s6-rc.d/nemo-asr-c/finish
#!/command/with-contenv bashio
# shellcheck shell=bash
# ============================================================================
# s6 finish for nemo-asr-c
# ============================================================================
bashio::log.info "nemo-asr-c stopped."
