# nemo-asr-c/rootfs/etc/s6-overlay/s6-rc.d/nemo-asr-c/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ============================================================================
# s6 longrun for nemo-asr-c — Wyoming STT server
# ============================================================================
set -euo pipefail

cd /usr/src/app || exit 1

args=(
    --uri "tcp://0.0.0.0:10370"
    --model "$(bashio::config 'model')"
    --quantization "$(bashio::config 'quantization')"
    --chunk-size "$(bashio::config 'chunk_size')"
)
bashio::config.true 'debug_logging' && args+=(--debug)

exec python3 -m wyoming_nemo_asr_c "${args[@]}"
