# zensical/rootfs/etc/s6-overlay/s6-rc.d/init-zensical/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Add-on: Zensical — initial seeding + first build.
#
# Seeds /config/docs/ from /opt/zensical/seed/ only when the directory is
# missing — the user's existing files are never overwritten. Then syncs
# into the working docs_dir and runs `zensical build --clean` so nginx
# has something to serve the moment it starts.
# ==============================================================================

SEED_DOCS=/opt/zensical/seed/docs
SEED_EXAMPLE_TOML=/opt/zensical/seed/zensical.example.toml
# First-run marker in the add-on-private /data (persists across restarts/
# updates). Once set, the example is never recreated — so deleting it makes
# it stay gone instead of reappearing every start.
SEED_MARKER=/data/.zensical-example-seeded
# Lives under /var/lib/zensical/ so its relative docs_dir/site_dir
# resolve to /var/lib/zensical/{docs,site}. Zensical 0.0.43 rejects
# absolute paths for those keys.
EFFECTIVE_TOML=/var/lib/zensical/zensical.toml

bashio::log.level "$(bashio::config 'log_level')"

if [ ! -d /config/docs ]; then
    mkdir -p /config/docs
    cp -r "$SEED_DOCS"/. /config/docs/ \
        || bashio::exit.nok 'Failed to seed /config/docs/'
    bashio::log.info 'Seeded /config/docs/ (no existing directory).'
fi

# Drop a publish-ready example into /config on the *first run only* (tracked
# by SEED_MARKER). The same file works both for the add-on and a GitHub Pages
# workflow. We only ever seed the example — never the active zensical.toml — so
# the bundled defaults stay in charge by default. After the first run we never
# recreate it, so a user who doesn't want it can delete it and it stays gone.
if [ ! -f "$SEED_MARKER" ] \
    && [ ! -f /config/zensical.example.toml ] \
    && [ ! -f /config/zensical.toml ]; then
    cp "$SEED_EXAMPLE_TOML" /config/zensical.example.toml
    bashio::log.info 'Wrote /config/zensical.example.toml — rename to zensical.toml to customise.'
fi
# Mark the first-run seed as done (even when skipped because the user already
# had a config), so it's a one-time event per install.
mkdir -p /data
: > "$SEED_MARKER"

CONFIG_FILE=$(/opt/zensical/bin/resolve-config.sh)
bashio::log.info "Using renderer config: $CONFIG_FILE"

/opt/zensical/bin/sync.sh \
    || bashio::exit.nok 'Initial sync failed'

# Rewrite docs_dir + site_dir to the container paths the add-on owns,
# preserving everything else from the source config. This is how a
# user's publish-ready /config/zensical.toml (with `docs_dir = "docs"`
# etc.) can render correctly inside the add-on.
/opt/zensical/bin/effective-config.py "$CONFIG_FILE" "$EFFECTIVE_TOML" \
    || bashio::exit.nok 'Failed to render effective config'

bashio::log.info 'Running initial zensical build…'
# --clean drops Zensical's incremental build cache so the cold start
# never reuses stale fragments from a previous container.
if ! zensical build --clean -f "$EFFECTIVE_TOML"; then
    bashio::log.warning 'Initial zensical build failed — nginx will serve a placeholder.'
    mkdir -p /var/lib/zensical/site
    printf '<!doctype html><meta charset=utf-8><title>Docs</title><h1>Build failed</h1><p>Check the add-on logs.</p>' \
        > /var/lib/zensical/site/index.html
fi
