# llm-conversation-agent/rootfs/etc/s6-overlay/s6-rc.d/llm-agent/run
#!/command/with-contenv bashio
# vim: ft=bash
# shellcheck shell=bash
# ==============================================================================
# Start the Wyoming LLM conversation agent.
# ==============================================================================

# Build CLI args from addon options. Optional fields are only forwarded when
# the user actually set them so the Python argparse defaults stay in charge.
flags=()
bashio::config.has_value 'request_timeout' \
    && flags+=('--request-timeout' "$(bashio::config 'request_timeout')")
bashio::config.has_value 'max_tool_iterations' \
    && flags+=('--max-tool-iterations' "$(bashio::config 'max_tool_iterations')")
bashio::config.has_value 'history_turns' \
    && flags+=('--history-turns' "$(bashio::config 'history_turns')")
# HA mcp_server delegation: only override the Python default when the
# user explicitly disabled it. Default (true) keeps the supervisor URL.
bashio::config.false 'ha_mcp_enabled' \
    && flags+=('--ha-mcp-url' 'disabled')
bashio::config.true 'debug_logging' && flags+=('--debug')

# TLS verification for outbound HTTP. Default true (only override when
# explicitly disabled).
bashio::config.false 'verify_ssl' \
    && flags+=('--verify-ssl' 'false')

# Embedding-based tool filtering. embedding_base_url / _api_key may be
# blank, in which case the Python side falls back to the chat values.
bashio::config.has_value 'embedding_base_url' \
    && flags+=('--embedding-base-url' "$(bashio::config 'embedding_base_url')")
bashio::config.has_value 'embedding_api_key' \
    && flags+=('--embedding-api-key' "$(bashio::config 'embedding_api_key')")
bashio::config.has_value 'embedding_model' \
    && flags+=('--embedding-model' "$(bashio::config 'embedding_model')")
bashio::config.has_value 'tool_filter_top_k' \
    && flags+=('--tool-filter-top-k' "$(bashio::config 'tool_filter_top_k')")
bashio::config.has_value 'ha_tools_mode' \
    && flags+=('--ha-tools-mode' "$(bashio::config 'ha_tools_mode')")

# Skills (agentskills.io format SKILL.md under /config/skills/<slug>/).
# /config/ is the addon_config mount, so /config/skills/ on the host is
# /addon_configs/{repo}_llm_conversation_agent/skills/. Pre-create the
# directory after the mount lands so users can drop files in without
# manual mkdir.
if bashio::config.true 'skills_enabled'; then
    flags+=('--skills-enabled' 'true')
    mkdir -p /config/skills
else
    flags+=('--skills-enabled' 'false')
fi

# Repeatable --skill-url flags. Each entry must be a GitHub repo URL
# (https://github.com/user/repo[@ref]) or a direct .tar.gz / .tgz / .zip
# archive URL. Hash-based update tracking lives at
# /data/skill_install_state.json (addon-private).
while read -r url; do
    [[ -n "${url}" ]] && flags+=('--skill-url' "${url}")
done < <(jq -r '.skill_urls[]?' /data/options.json)

# Per-user persistent memory (R2). Mirrors the skills directory
# convention: /config/memory/ on the host = /addon_configs/{repo}_llm_
# conversation_agent/memory/. Pre-create so the first memory_save
# never races a missing parent.
if bashio::config.true 'memory_enabled'; then
    flags+=('--memory-enabled' 'true')
    mkdir -p /config/memory
else
    flags+=('--memory-enabled' 'false')
fi

# Workspace (R3) — SOUL.md / IDENTITY.md / HEARTBEAT.md / BOOTSTRAP.md
# at /config/ root. /config is the supervisor-mounted addon_config dir
# so it always exists; workspace is always on (no toggle).

# Prometheus /metrics endpoint.
if bashio::config.true 'metrics_enabled'; then
    flags+=('--metrics-enabled' 'true')
else
    flags+=('--metrics-enabled' 'false')
fi
bashio::config.has_value 'metrics_port' \
    && flags+=('--metrics-port' "$(bashio::config 'metrics_port')")

# Repeatable --language flags.
if bashio::config.has_value 'languages'; then
    while read -r lang; do
        [[ -n "${lang}" ]] && flags+=('--language' "${lang}")
    done < <(bashio::config 'languages|keys' \
                || jq -r '.languages[]?' /data/options.json)
fi

# Repeatable --extra-header NAME=VALUE flags.
while read -r header; do
    [[ -n "${header}" ]] && flags+=('--extra-header' "${header}")
done < <(jq -r '(.extra_headers // [])[] | "\(.name)=\(.value)"' /data/options.json)

bashio::log.info "Starting Wyoming LLM conversation agent"
bashio::log.info "  upstream:  $(bashio::config 'base_url')"
bashio::log.info "  model:     $(bashio::config 'model')"
if bashio::config.true 'ha_mcp_enabled'; then
    bashio::log.info "  HA MCP:    enabled (supervisor proxy)"
else
    bashio::log.info "  HA MCP:    disabled"
fi

# Extract mcp_servers as JSON file (preserves nested structure that
# bashio can't pass through CLI args without mangling).
mkdir -p /tmp/llm-agent
jq -c '.mcp_servers // []' /data/options.json > /tmp/llm-agent/mcp_servers.json

exec python3 -m wyoming_llm_agent \
    --uri 'tcp://0.0.0.0:10500' \
    --base-url "$(bashio::config 'base_url')" \
    --api-key "$(bashio::config 'api_key')" \
    --model "$(bashio::config 'model')" \
    --mcp-servers '@/tmp/llm-agent/mcp_servers.json' \
    "${flags[@]}"
