# rethink/rootfs/etc/s6-overlay/s6-rc.d/rethink/run
#!/command/with-contenv bashio
# vim: ft=bash
# shellcheck shell=bash
# ==============================================================================
# Render rethink config.json from add-on options + MQTT credentials, then
# launch rethink-cloud. Re-runs on every container start, so editing
# /data/config.json by hand has no effect — change add-on options instead.
# ==============================================================================

if ! bashio::services.available "mqtt"; then
    bashio::exit.nok "MQTT service unavailable. Install and start the Mosquitto broker add-on first."
fi

MQTT_HOST=$(bashio::services mqtt "host")
MQTT_PORT=$(bashio::services mqtt "port")
MQTT_USER=$(bashio::services mqtt "username")
MQTT_PASS=$(bashio::services mqtt "password")

mkdir -p /data/state

jq -n \
    --slurpfile opts /data/options.json \
    --arg mqtt_url  "mqtt://${MQTT_HOST}:${MQTT_PORT}" \
    --arg mqtt_user "${MQTT_USER}" \
    --arg mqtt_pass "${MQTT_PASS}" \
    '{
        hostname: $opts[0].hostname,
        homeassistant: {
            mqtt_url:          $mqtt_url,
            discovery_prefix:  $opts[0].discovery_prefix,
            rethink_prefix:    $opts[0].rethink_prefix,
            mqtt_user:         $mqtt_user,
            mqtt_pass:         $mqtt_pass
        },
        ca_key_file:       "/data/ca.key",
        ca_cert_file:      "/data/ca.cert",
        https_port:        $opts[0].https_port,
        mqtts_port:        $opts[0].mqtts_port,
        mqtt_port:         $opts[0].mqtt_port,
        thinq1_https_port: $opts[0].thinq1_https_port,
        thinq1_port:       $opts[0].thinq1_port,
        management_port:   $opts[0].management_port,
        bridge: { storage_path: "/data/state" },
        log:    $opts[0].log_levels
    }' > /data/config.json

MGMT_PORT=$(jq -r .management_port /data/config.json)
bashio::log.info "hostname        : $(jq -r .hostname /data/config.json)"
bashio::log.info "MQTT broker     : ${MQTT_HOST}:${MQTT_PORT} (user: ${MQTT_USER})"
bashio::log.info "management UI   : http://<HA-IP>:${MGMT_PORT}  (NO AUTH — restrict on LAN)"

cd /app/rethink
exec node dist/rethink-cloud.js /data/config.json
