# wardrowbe/rootfs/etc/s6-overlay/s6-rc.d/postgres/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash
set -euo pipefail
mkdir -p /run/postgresql
exec env LD_PRELOAD=/usr/local/lib/libfakeeuid.so \
  postgres -D /var/lib/postgresql/data
