# zensical/rootfs/etc/s6-overlay/s6-rc.d/zensical-watcher/finish
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# zensical-watcher finish handler.
#
# s6-supervise restarts the service ONLY when finish exits 0.
# - Signal kill / exit 0 → halt the container (manual stop or graceful shutdown).
# - Non-zero exit → exit 0 (restart the watcher).  The watcher is stateless
#   (inotifywait + zensical build); a crash is almost always a transient resource
#   issue.  We don't halt the container so nginx keeps serving the last built
#   site and the user can read logs.
# ==============================================================================
# shellcheck disable=SC2155
readonly exit_code_container=$(</run/s6-linux-init-container-results/exitcode)
readonly exit_code_service="${1}"
readonly exit_code_signal="${2}"

bashio::log.info \
  "zensical-watcher exited with code ${exit_code_service}" \
  "(by signal ${exit_code_signal})"

if [[ "${exit_code_signal}" -ne 0 ]]; then
  # Killed by signal — transient (manual stop / resource limit).  Halt.
  if [[ "${exit_code_container}" -eq 0 ]]; then
    echo $((128 + exit_code_signal)) > /run/s6-linux-init-container-results/exitcode
  fi
  exec /run/s6/basedir/bin/halt
elif [[ "${exit_code_service}" -eq 0 ]]; then
  # Clean exit — transient (graceful shutdown).  Halt.
  exec /run/s6/basedir/bin/halt
else
  # Non-zero exit — restart the watcher via exit 0.  Stateless, so restart
  # is safe; a config error would cause an obvious crash-loop.
  bashio::log.warning \
    "zensical-watcher: crashed with exit ${exit_code_service} — restarting"
  exit 0
fi
