# zensical/rootfs/etc/s6-overlay/s6-rc.d/zensical-watcher/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Watch /config/docs/ and /addon_config/ for changes; re-sync into the
# renderer docs dir and rebuild the site. Debounced to one rebuild per 1s
# burst of inotify events.
#
# Pattern is modelled on wendevlin/homedocs (chokidar watcher + per-change
# rebuild), adapted to Zensical.
# ==============================================================================

RESOLVE=/opt/zensical/bin/resolve-config.sh
EFFECTIVE=/opt/zensical/bin/effective-config.py
SYNC=/opt/zensical/bin/sync.sh
EFFECTIVE_TOML=/var/lib/zensical/zensical.toml
DEBOUNCE=1

bashio::log.level "$(bashio::config 'log_level')"

build() {
    # Re-create the watched paths if the user deleted them — otherwise
    # inotifywait exits with code 1 and s6 restarts us in a tight loop.
    mkdir -p /config/docs /addon_config

    "$SYNC" || {
        bashio::log.warning 'sync failed; skipping this rebuild'
        return 0
    }

    # Resolve + rewrite on every build so that creating/removing
    # /config/zensical.toml or /addon_config/zensical.toml takes effect
    # immediately, and the user's docs_dir/site_dir get replaced with
    # the container paths nginx and sync.sh rely on.
    local config_file
    config_file=$("$RESOLVE")
    "$EFFECTIVE" "$config_file" "$EFFECTIVE_TOML" || {
        bashio::log.warning 'effective config rewrite failed; skipping this rebuild'
        return 0
    }

    # No --clean here: warm builds reuse Zensical's incremental cache for
    # millisecond rebuilds. init-zensical handles the cold-start clean.
    if zensical build -f "$EFFECTIVE_TOML"; then
        bashio::log.info "site rebuilt (source: $config_file)"
    else
        bashio::log.warning 'zensical build failed (previous site still served)'
    fi
}

# Build the inotify watch list. /config/zensical.toml is only included
# if it exists right now — inotifywait fails the whole invocation on
# missing paths. The list is rebuilt every iteration, so the moment a
# /config/docs edit fires the watcher, a *newly* created
# /config/zensical.toml gets picked up on the next pass.
collect_watch_paths() {
    local p=(/config/docs /addon_config)
    [ -f /config/zensical.toml ] && p+=(/config/zensical.toml)
    printf '%s\n' "${p[@]}"
}

bashio::log.info "watching /config/docs/[,zensical.toml] and /addon_config/ (debounce ${DEBOUNCE}s)"

while true; do
    mapfile -t WATCH_PATHS < <(collect_watch_paths)

    # Block until first event.
    inotifywait -r -q \
        -e modify,create,delete,move,moved_to,moved_from \
        "${WATCH_PATHS[@]}" >/dev/null 2>&1 || true

    # Drain follow-up events for $DEBOUNCE seconds — collapses bursts
    # (multi-file saves, editor swap files) into a single rebuild.
    while inotifywait -r -q -t "$DEBOUNCE" \
        -e modify,create,delete,move,moved_to,moved_from \
        "${WATCH_PATHS[@]}" >/dev/null 2>&1; do :; done

    build
done
