# wardrowbe/rootfs/etc/s6-overlay/s6-rc.d/frontend/finish
#!/usr/bin/with-contenv bash
exitcode="${1:-0}"
signal="${2:-0}"
echo "[frontend] stopped (exit=${exitcode}, signal=${signal})"
exit 0
