# wardrowbe/rootfs/etc/s6-overlay/s6-rc.d/backend/finish
#!/usr/bin/with-contenv bash
# shellcheck shell=bash
exitcode="${1:-0}"
signal="${2:-0}"
echo "[backend] stopped (exit=${exitcode}, signal=${signal})"
exit 0
