# nemo-asr-cpp/rootfs/etc/s6-overlay/s6-rc.d/nemo-asr-cpp/run
#!/command/with-contenv bashio
# vim: ft=bash
# shellcheck shell=bash
# ==============================================================================
# Start the NeMo ASR (cpp) Wyoming service — parakeet.cpp (ggml) via ctypes.
# ==============================================================================

# GGUF model downloads to the persistent /data volume (not the image).
export HOME=/data
mkdir -p /data/models

flags=()
bashio::config.has_value 'model'         && flags+=('--model'         "$(bashio::config 'model')")
bashio::config.has_value 'quantization'  && flags+=('--quantization'  "$(bashio::config 'quantization')")
bashio::config.has_value 'chunk_size'    && flags+=('--chunk-size'    "$(bashio::config 'chunk_size')")
bashio::config.has_value 'hotword_boost' && flags+=('--hotword-boost' "$(bashio::config 'hotword_boost')")
bashio::config.has_value 'hf_token'      && flags+=('--hf-token'      "$(bashio::config 'hf_token')")
bashio::config.true 'debug_logging'      && flags+=('--debug')

# Hotwords are a YAML list; join into one newline-separated --hotwords value.
hotwords=$(jq -r '(.hotwords // []) | join("\n")' /data/options.json 2>/dev/null || echo "")
[ -n "${hotwords}" ] && flags+=('--hotwords' "${hotwords}")

# Quote the array so flag values containing spaces (e.g. "Bahasa Indonesia",
# multi-line --hotwords) stay one argument each.
exec python3 -m wyoming_nemo_asr_cpp \
    --uri 'tcp://0.0.0.0:10360' \
    "${flags[@]}"
